# assets/etc/supervisor/conf.d/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
stdout_logfile=NONE
